# sv_snap.conf
[program:pi-snap]
command = bash /home/pi/snap/sshreverse.sh
user = pi
autostart=true
autorestart=true
